# System_Config/healthcheck.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL__</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__WORKSPACE_ROOT__/System_Config/healthcheck.sh</string>
    </array>

    <!-- Regenerate the status page at login and every 4 hours thereafter, so an
         open browser tab (which auto-refreshes every 5 min) always reflects a
         recent probe. The check is read-only + cheap; running it often is safe. -->
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>14400</integer>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/healthcheck.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/healthcheck.launchd.err</string>
</dict>
</plist>
