# System_Config/mondayinit.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL__</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__WORKSPACE_ROOT__/System_Config/monday_init.sh</string>
    </array>

    <!-- Runs at login/startup (RunAtLoad) AND every Monday 08:00 (Weekday 1) so the
         week's note exists whether the Mac was off over the weekend or already awake.
         monday_init.sh is idempotent — it skips if this ISO week's note already
         exists — so the login-time re-run is safe. -->
    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/monday_init.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/monday_init.launchd.err</string>
</dict>
</plist>
